# systemd/canlogger.service
[Unit]
Description=CANlogger waveform capture service
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/CANlogger
ExecStart=/usr/bin/env bash -lc 'docker compose -f /opt/CANlogger/docker-compose.yml up -d'
ExecStop=/usr/bin/env bash -lc 'docker compose -f /opt/CANlogger/docker-compose.yml down'
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
